# apps/trace-parser/com.paperclip.trace-parser.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.paperclip.trace-parser</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/michaeldavidson/Developer/paperclip/apps/trace-parser/bin/run-parse.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>60</integer>
    <key>StandardOutPath</key>
    <string>/tmp/paperclip-trace-parser.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/paperclip-trace-parser.stderr.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
